# linux/rke2-before-shutdown-service/rke2-pre-shutdown.service
[Unit]
Description=Pre-shutdown hook to update NamespaceLifecyclePolicy CRs
DefaultDependencies=no

# reversed logic: systemd stops things in REVERSE order of 'After'
# This ensures rke2-server is still running when this service stops
After=rke2-server.service network-online.target
Before=shutdown.target reboot.target halt.target

[Service]
Type=oneshot
# ExecStart does nothing, just marks the service as "active"
ExecStart=/bin/true

# This is where your actual script runs
ExecStop=/usr/local/bin/rke2-pre-shutdown.sh

# This is critical: keeps the service "active" so ExecStop can trigger later
RemainAfterExit=yes

# Give your script enough time to scale down all services
TimeoutStopSec=600 

[Install]
# Start this during normal boot so it's ready to "stop" during shutdown
WantedBy=multi-user.target